# umbra.service
[Unit]
Description=Umbra shadow engine — CDP server
Documentation=https://github.com/Celebez/umbra
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Path to the Umbra engine binary (the engine core) and the umbra venv/python.
Environment=UMBRA_ENGINE_BIN=/usr/local/bin/umbra-engine
Environment=PATH=/opt/umbra/.venv/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/opt/umbra
# Binds 127.0.0.1:9222 (loopback only — NOT exposed on a public IP/port).
ExecStart=/opt/umbra/.venv/bin/umbra serve --port 9222 --stealth
Restart=on-failure
RestartSec=3
# Hardened: no new privileges, private /tmp.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/umbra
PrivateTmp=true

[Install]
WantedBy=multi-user.target
